# .build/builder/iva/supervisord.conf
; RUN COMMAND
; supervisor -c /path/to/supervisord.conf

[unix_http_server]
file=/opt/supervisor_http.sock  ; path to your socket file

[supervisord]
logfile=/opt/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB           ; maximum size of logfile before rotation
logfile_backups=10              ; number of backed up logfiles
loglevel=error                  ; info, debug, warn, trace
pidfile=/opt/supervisord.pid    ; pidfile location
nodaemon=true                   ; run supervisord as a daemon
minfds=1024                     ; number of startup file descriptors
minprocs=200                    ; number of process descriptors
user=root                       ; default user

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///opt/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:video-pipeline]
command=/bin/bash -c "/src/build/gstreamer"
priority=5
autostart=true
startsecs=1
startretries=1000
autorestart=true
stopwaitseconds=5
; stdout_logfile=/src/logs/stdout.log
; stderr_logfile=/src/logs/stderr.log
; logfile_maxbytes=50MB
; stderr_logfile_maxbytes=50MB

stdout_logfile=/dev/stdout   ; when directing to /dev/stdout, must have maxbytes=0
stderr_logfile=/dev/stderr
logfile_maxbytes=0
stderr_logfile_maxbytes = 0



